# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:postgresql]
command=/usr/bin/postgres -D %(ENV_PGDATA)s
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql-error.log

[program:tr-engine]
command=/app/tr-engine --config /app/config/config.yaml
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=3
priority=20
stdout_logfile=/var/log/supervisor/tr-engine.log
stderr_logfile=/var/log/supervisor/tr-engine-error.log
environment=DB_HOST="localhost",DB_PORT="5432",DB_USER="%(ENV_POSTGRES_USER)s",DB_PASSWORD="%(ENV_POSTGRES_PASSWORD)s",DB_NAME="%(ENV_POSTGRES_DB)s"

[eventlistener:postgresql_ready]
command=sh -c "sleep 5"
events=PROCESS_STATE
